# deploy/systemd/dev/wolf-dashboard.service
[Unit]
Description=Wolf dashboard — the edge component (dev, user-level)
Documentation=file://@REPO_ROOT@/docs/decisions/0016-wolf-component-architecture-and-packaging.md
# Per ADR 0016 v3: fully independent — no After=/Requires=/Wants= between Wolf units.
# The dashboard's reverse proxy (Phase 5.6-a) auto-retries outbound
# fetches to wolf-server; nothing breaks if wolf-server isn't up yet.

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@/services/dashboard

# Dev mode uses `npm run dev`, which invokes scripts/dev.mjs and
# picks up the HTTPS + mTLS auto-detect logic from Phase 5.4-d /
# 5.6-c. The install-user-systemd Makefile target substitutes
# @NODE_BIN@ with the directory containing the operator's `npm`
# (typically nvm: `~/.nvm/versions/node/<v>/bin`); we extend PATH
# so npm + node + npx + every other Node tool resolves cleanly.
Environment=PATH=@NODE_BIN@:/usr/bin:/usr/local/bin:/bin
ExecStart=@NODE_BIN@/npm run dev

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
